# qemu-scripts.service
[Unit]
Description=Qemu-Scripts Server

[Service]
Type=simple
ExecStart=/opt/qemu-scripts/start.sh

[Install]
WantedBy=multi-user.target
